# s6-overlay/s6-rc.d/init-adguardhome/run
#!/command/with-contenv sh
exec /app/AdGuardHome -c /var/lib/app/adguardhome/AdGuardHome.yaml -w /var/lib/app/adguardhome | sed 's/^/[adguardhome] /'