# s6-overlay/s6-rc.d/init-alertmanager/run
#!/command/with-contenv sh
exec /usr/bin/alertmanager --config.file=/app/monitor/alertmanager-config.yml | sed 's/^/[alertmanager] /'