# s6-overlay/s6-rc.d/init-prometheus-node-exporter/run
#!/command/with-contenv sh
exec /usr/bin/node_exporter --web.listen-address=127.0.0.1:9100 | sed 's/^/[node_exporter] /'