# s6-overlay/s6-rc.d/init-pushgateway/run
#!/command/with-contenv sh
exec /usr/bin/pushgateway --web.listen-address=127.0.0.1:9091 | sed 's/^/[pushgateway] /'