# s6-overlay/s6-rc.d/init-promtail/run
#!/command/with-contenv sh
exec /usr/bin/promtail -config.file=/app/monitor/promtail-config.yml