# s6-overlay/s6-rc.d/init-caddy/run
#!/command/with-contenv sh
exec caddy run --config /etc/caddy/Caddyfile