# s6-overlay/s6-rc.d/init-tailscale/run
#!/command/with-contenv sh
/app/tailscaled --state=/var/lib/app/tailscale/tailscaled.state --socket=/var/run/tailscale/tailscaled.sock | sed 's/^/[tailscale] /'