# s6-overlay/s6-rc.d/init-tailscale-up/run
#!/command/with-contenv sh
/app/tailscale up --authkey=${TAILSCALE_AUTHKEY} --hostname=internal --ssh=true --advertise-exit-node=true

tail -f /dev/null