# s6-overlay/s6-rc.d/init-loki/run
#!/command/with-contenv sh
exec /usr/bin/loki -config.file=/app/monitor/loki-config.yml | sed 's/^/[loki] /'