# s6-overlay/s6-rc.d/init-prometheus/run
#!/command/with-contenv sh
exec /usr/bin/prometheus --config.file=/app/monitor/prometheus-config.yml