# s6-overlay/s6-rc.d/init-adguard-exporter/run
#!/command/with-contenv sh
exec /usr/bin/adguard-exporter