# s6-overlay/s6-rc.d/init-grafana/run
#!/command/with-contenv sh
exec /var/lib/app/grafana/bin/grafana server --homepath /var/lib/app/grafana | sed 's/^/[grafana] /'